# etc/systemd/system/undervolt.service
[Unit]
Description=undervolt

[Service]
Type=oneshot
ExecStart=/usr/local/bin/undervolt -v --core -120 --cache -120
